# recipes/redis/supervisord-redis.conf
[program:redis]
command=/bin/sh $PREFIX/etc/supervisord/startup/redis

startsecs = 5
pidfile = $PREFIX/var/run/redis/redis.pid
autostart = true
autorestart = unexpected
stopasgroup = false
stopsignal = TERM

stdout_logfile = $PREFIX/var/log/redis/redis.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 5
redirect_stderr = true
